# sv/zramen/run
#!/bin/sh
[ -r ./conf ] && . ./conf
zramen make
exec pause
